# HomeAppDocker/webservers/php-cli/cron/crontab
# crontab to run command every minute
* * * * * /var/www/html/bin/console app:trigger-check