# HomeAppDocker/apache/supervisor/update-current-reading.ini
[program:messenger-consume-current-reading-update]
command=php bin/console rabbitmq:consumer upload_current_reading_sensor_data -m 200
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d